# cron_daily_report.cron
# Update /path/to/repo with the absolute path to this repository.
# Optional: set timezone for the 20:00 run (e.g. CRON_TZ=Europe/Paris).
# CRON_TZ=Europe/Paris
0 20 * * * cd /path/to/repo/project && /usr/bin/env python3 report.py --ticker AAPL --interval 1d --output-dir /path/to/repo/project/reports >> /path/to/repo/project/logs/daily_report.log 2>&1
